# fcpi.service
[Unit]
Description=fcpi

[Service]
Type=simple
ExecStart=/home/pi/fcpi/fcpi.sh
Restart=always

[Install]
WantedBy=multi-user.target
